# meta-wpe-distro/recipes-graphics/wayland/files/wayland-1.path
[Unit] 
Description=Check wayland-1 socket exists
RequiresMountsFor=/run/user/1000
Requires=weston.service
After=weston.service

[Path]
PathExists=/run/user/1000/wayland-1

[Install]
WantedBy=graphical.target

